# docs/configuracion/crontab/crontab
# m h  dom mon dow   command
* * * * * perl /usr/local/koha/intranet/scripts/z3950daemon/procesarColaZ3950.pl &
0 3 * * 0 perl /usr/local/koha/intranet/scripts/obtener_portadas_de_registros.pl &
